# hetnet/neo4j/deployment/neo4j-docker.conf
[program:neo4j-docker]
command=bash neo4j-supervisor.sh
directory=/home/ubuntu/docker-scripts/
user=ubuntu
group=ubuntu
autostart=true
autorestart=true
priority=991
stopsignal=KILL
